# data/helexa-router.service
[Unit]
Description=helexa-router — federation data plane (capacity-aware dispatch across operator cortexes)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/helexa-router serve --config /etc/helexa-router/helexa-router.toml
# Stateless streaming proxy: polls each configured cortex's /health +
# /v1/models topology and dispatches inference on capacity. Restart
# unconditionally if it ever exits.
Restart=always
RestartSec=10
User=helexa-router
Group=helexa-router
StateDirectory=helexa-router
StateDirectoryMode=0755

[Install]
WantedBy=multi-user.target
